# claude-code-terminal/rootfs/etc/services.d/ttyd/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "=== Claude Code Terminal starting ==="

# -------------------------------------------------------
# Persistent storage setup
# -------------------------------------------------------
PERSIST_DIR=/data/home/.claude
mkdir -p "${PERSIST_DIR}" /data/home/.config/claude-code /data/projects

# -------------------------------------------------------
# Import credentials from shared config if available
# -------------------------------------------------------
if [ -f /homeassistant/.claude-creds-transfer.json ]; then
    bashio::log.info "Found credential transfer file - importing..."
    cp /homeassistant/.claude-creds-transfer.json "${PERSIST_DIR}/.credentials.json"
    chmod 600 "${PERSIST_DIR}/.credentials.json"
    rm -f /homeassistant/.claude-creds-transfer.json
    bashio::log.info "Credentials imported successfully"
fi

# -------------------------------------------------------
# Read add-on options
# -------------------------------------------------------
API_KEY=$(bashio::config 'anthropic_api_key' '')
MODEL=$(bashio::config 'model' 'sonnet')
YOLO_MODE=$(bashio::config 'yolo_mode' 'false')
SESSION_NAME=$(bashio::config 'session_name' 'claude')

# -------------------------------------------------------
# API key (if provided, overrides OAuth credentials)
# -------------------------------------------------------
if bashio::var.has_value "${API_KEY}"; then
    export ANTHROPIC_API_KEY="${API_KEY}"
    bashio::log.info "Anthropic API key configured"
else
    bashio::log.info "No API key set - using OAuth credentials"
fi

# -------------------------------------------------------
# Symlink Claude config to persistent storage
# -------------------------------------------------------
if [ ! -L /root/.claude ]; then
    rm -rf /root/.claude
    ln -s "${PERSIST_DIR}" /root/.claude
fi
mkdir -p /root/.config
if [ ! -L /root/.config/claude-code ]; then
    rm -rf /root/.config/claude-code
    ln -s /data/home/.config/claude-code /root/.config/claude-code
fi

# -------------------------------------------------------
# YOLO mode
# -------------------------------------------------------
if bashio::var.true "${YOLO_MODE}"; then
    bashio::log.info "YOLO mode enabled - all permissions auto-approved"
    echo '{"permissions":{"allow":["Bash(*)","Read(*)","Write(*)","Edit(*)","Glob(*)","Grep(*)"],"deny":[]}}' \
        > "${PERSIST_DIR}/settings.json"
fi

# -------------------------------------------------------
# Environment file for bash sessions
# -------------------------------------------------------
ENV_FILE=/data/env.sh
cat > "${ENV_FILE}" <<ENVEOF
export TERM=xterm-256color
export LANG=C.UTF-8
export SHELL=/bin/bash
export HOME=/root
export CLAUDE_MODEL="${MODEL}"
export PATH="/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/root/.npm-global/bin:\${PATH}"
ENVEOF

if bashio::var.has_value "${API_KEY}"; then
    echo "export ANTHROPIC_API_KEY=\"${API_KEY}\"" >> "${ENV_FILE}"
fi

if [ -n "${SUPERVISOR_TOKEN:-}" ]; then
    echo "export SUPERVISOR_TOKEN=\"${SUPERVISOR_TOKEN}\"" >> "${ENV_FILE}"
fi

# -------------------------------------------------------
# Register HA MCP server with Claude Code
# -------------------------------------------------------
if bashio::var.has_value "${SUPERVISOR_TOKEN:-}"; then
    bashio::log.info "Configuring HA MCP server..."
    export SUPERVISOR_TOKEN="${SUPERVISOR_TOKEN}"
    claude mcp remove homeassistant -s user 2>/dev/null || true
    claude mcp add homeassistant -s user \
        -- node /opt/ha-mcp-server/build/index.js 2>/dev/null \
        && bashio::log.info "HA MCP server registered (13 tools)" \
        || bashio::log.warning "HA MCP registration failed (non-fatal)"
fi

# -------------------------------------------------------
# Set onboarding complete flag (prevents first-run wizard)
# Must run AFTER mcp config since claude mcp add writes to .claude.json
# -------------------------------------------------------
if [ -f /root/.claude.json ]; then
    TMP=$(jq '. + {"hasCompletedOnboarding":true}' /root/.claude.json 2>/dev/null) && echo "${TMP}" > /root/.claude.json
else
    echo '{"hasCompletedOnboarding":true}' > /root/.claude.json
fi
bashio::log.info "Onboarding flag set"

# -------------------------------------------------------
# Run auto-discovery
# -------------------------------------------------------
if [ -n "${SUPERVISOR_TOKEN:-}" ]; then
    bashio::log.info "Running HA auto-discovery..."
    /opt/scripts/ha-discovery.sh 2>/dev/null \
        && bashio::log.info "HA context written to /homeassistant/.claude/ha_context.md" \
        || bashio::log.warning "Auto-discovery failed (non-fatal)"
fi

# -------------------------------------------------------
# Deploy CLAUDE.md to projects dir if not present
# -------------------------------------------------------
if [ ! -f /data/projects/CLAUDE.md ]; then
    cp /opt/CLAUDE.md /data/projects/CLAUDE.md
    bashio::log.info "Deployed CLAUDE.md template"
fi

# -------------------------------------------------------
# Shell configuration
# -------------------------------------------------------
export SHELL=/bin/bash
echo "set-option -g default-shell /bin/bash" >> /root/.tmux.conf
echo "set-option -g default-command /bin/bash" >> /root/.tmux.conf

cat > /root/.bashrc <<'RCEOF'
[ -f /data/env.sh ] && source /data/env.sh
alias c="claude"
alias cc="claude --continue"
alias cr="claude --resume"

# Welcome message
echo ""
echo "  Claude Code Terminal for Home Assistant"
echo "  ─────────────────────────────────────────"
echo "  Commands:  c (claude)  |  cc (continue)  |  cr (resume)"
echo "  Docs:      /homeassistant/.claude/"
echo "  Projects:  /data/projects/"
echo ""
RCEOF

# -------------------------------------------------------
# Launch ttyd → tmux → bash
# -------------------------------------------------------
ingress_port=$(bashio::addon.ingress_port)
bashio::log.info "Starting ttyd on ingress port ${ingress_port}"

cd /data/projects || cd /root

exec ttyd \
    -i 0.0.0.0 \
    -p "${ingress_port}" \
    -W \
    -t fontSize=14 \
    -t fontFamily="Monaco,Consolas,monospace" \
    /bin/bash -c "exec tmux new-session -A -s ${SESSION_NAME}"
